# ai-news-dashboard.service
[Unit]
Description=AI News Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dave/ai-news-dashboard
ExecStart=/usr/bin/node /home/dave/ai-news-dashboard/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
